# scripts/systemd/current-link-pc-mirror-sync.service
[Unit]
Description=Current Link ERP PC mirror sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/current-link/app
ExecStart=/usr/bin/python /opt/current-link/app/scripts/run_pc_mirror_sync.py
